# systemd/tms.service
# ====================================
# Systemd Service for TMS v2.0
# ====================================
#
# Installation:
#   sudo cp systemd/tms.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable tms
#   sudo systemctl start tms
#
# Management:
#   sudo systemctl status tms
#   sudo systemctl restart tms
#   sudo systemctl stop tms
#   sudo journalctl -u tms -f

[Unit]
Description=TMS v2.0 - Oracle Tuning Management System
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/tms
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/var/www/tms/.env.production
ExecStart=/usr/bin/node /var/www/tms/node_modules/next/dist/bin/next start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tms

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/tms/logs

# Resource Limits
LimitNOFILE=65536
MemoryLimit=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
